# incubator/scripts/com.tokentank.morning-tweet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tokentank.morning-tweet</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/bart/Documents/code/vibeceo/sms-bot && source .env.local && export TOKENTANK_AUTO_TWEET=1 && /opt/homebrew/bin/npx tsx ../incubator/scripts/auto-tweet.ts morning</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/bart/Documents/code/vibeceo/incubator/scripts/morning-tweet.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bart/Documents/code/vibeceo/incubator/scripts/morning-tweet.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/bart/Documents/code/vibeceo</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
